# systemd/update-unbound-blacklist.service
[Unit]
Description=Update dns blacklist used by unbound adblocker

[Service]
Type=simple
ExecStart=/usr/local/bin/update-unbound-blacklist.sh
